# meta-jetson/recipes-software/maxpower/maxpower/maxpower.init
#! /bin/sh
### BEGIN INIT INFO
# Provides:    maxpower
# Description: Maximum GPU performance setup
### END INIT INFO

case "$1" in
    start)
        /usr/sbin/nvpmodel -m 0 -f /etc/nvpmodel.conf
        /usr/bin/jetson_clocks
        /usr/bin/jetson_clocks --fan
        ;;
    stop)
        ;;
    *)
        echo "Usage: $0 (start|stop)"
        ;;
esac

exit 0
