# meta-rdfm/recipes-rdfm/rdfm-client/files/rdfm-daemon.service
[Unit]
Description=RDFM's Linux client daemon
After=network.target

[Service]
Type=simple
ExecStart=rdfm daemonize
